# range_libc/pywrapper/ros/build/devel/etc/catkin/profile.d/1.ros_etc_dir.sh
# generated from ros_environment/env-hooks/1.ros_etc_dir.sh.em

# env variable in develspace
export ROS_ETC_DIR="/Users/anna/Documents/TUM/Thesis/ba-bio-inspired-navigation-main/range_libc/pywrapper/ros/build/devel/etc/ros"
